# templates/acme-sh.timer
[Unit]
Description=Schedules acme.sh.service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
